# debian/postinst
#!/bin/sh

set -e

case "$1" in
    configure)
        
	chmod +x /usr/sbin/so-apt-check /usr/sbin/sostat* /usr/sbin/soup /usr/sbin/so-purge-old-kernels || echo "Unable to set sostat utilities executable."
	
	FILE="/etc/update-manager/release-upgrades"
	if [ -f $FILE ]; then
		if grep "^Prompt=never" $FILE >/dev/null 2>&1; then
			echo "$FILE already configured to avoid prompting for next LTS release."
		else
			echo "Updating $FILE to avoid prompting for next LTS release."
			sed -i 's|^Prompt=.*$|Prompt=never|g' $FILE || echo "Error updating $FILE"
		fi
	fi

	FILE="/var/ossec/etc/ossec.conf"
	if grep "/usr/bin/sostat-interface" $FILE >/dev/null 2>&1; then
		echo "Updating sostat-interface path in $FILE."
		sed -i 's|/usr/bin/sostat-interface|/usr/sbin/sostat-interface|g' $FILE || echo "Error updating $FILE."
	fi
        
	# Emergency patch for redis-server issue
	# https://github.com/Security-Onion-Solutions/security-onion/issues/1386
	if grep "LOGSTASH_OUTPUT_REDIS=yes" /etc/nsm/securityonion.conf >/dev/null 2>&1; then
		if ! pgrep redis-server >/dev/null 2>&1; then
			echo "LOGSTASH_OUTPUT_REDIS=yes but redis-server doesn't appear to be running."
			echo "Enabling and starting redis-server."
			update-rc.d -f redis-server enable || echo "Error enabling redis-server."
			service redis-server start || echo "Error starting redis-server."
		fi
	fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;


    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
